# project_backups/backup_20251009_002504/backup_daemon.service
[Unit]
Description=Project Auto-Backup Daemon
After=network.target

[Service]
Type=forking
User=%i
WorkingDirectory=/mnt/c/webui
ExecStart=/mnt/c/webui/auto_backup_daemon.sh
ExecStop=/mnt/c/webui/stop_backup_daemon.sh
Restart=on-failure
RestartSec=10
StandardOutput=null
StandardError=null

# Resource limits for optimization
MemoryMax=256M
CPUQuota=20%
Nice=19
IOSchedulingClass=idle

[Install]
WantedBy=multi-user.target

